# server_configs/gunicorn_codomodo.socket
[Unit]
Description=gunicorn socket for kodjaz project

[Socket]
ListenStream=/run/gunicorn/kodjaz_socket

[Install]
WantedBy=sockets.target
