# server_configs/gunicorn_codomodo.socket
[Unit]
Description=gunicorn socket for codomodo project

[Socket]
ListenStream=/run/gunicorn/codomodo_socket

[Install]
WantedBy=sockets.target
