# server_configs/gunicorn_kodjaz.socket
[Unit]
Description=gunicorn socket for kodjaz project

[Socket]
ListenStream=/run/gunicorn/example_sock

[Install]
WantedBy=sockets.target
